# private_dot_config/systemd/user/jenv-refresh.path
[Unit]
Description=Detects when installed jdk's change and triggers jenv to update

[Path]
PathChanged=/usr/lib/jvm
Unit=jenv-refresh.service

[Install]
WantedBy=default.target
